# root/etc/services.d/aria2/run
#!/usr/bin/with-contenv bash

mkdir -p /home/app/download
chown -R app:app /home/app

# Create aria2.conf
if [ $aria2Secret ]; then
	echo "rpc-secret=${aria2Secret}" >> /etc/aria2/aria2.conf
fi
if [ $aria2MaxConcurrentDownloads ]; then
	echo "max-concurrent-downloads=${aria2MaxConcurrentDownloads}" >> /etc/aria2/aria2.conf
fi
if [ $aria2SeedRatio ]; then
	echo "seed-ratio=${aria2SeedRatio}" >> /etc/aria2/aria2.conf
fi
if [ $aria2SeedTime ]; then
	echo "seed-time=${aria2SeedTime}" >> /etc/aria2/aria2.conf
fi
if [ $aria2DisableIpv6 ]; then
	echo "disable-ipv6=${aria2DisableIpv6}" >> /etc/aria2/aria2.conf
fi
if [ $aria2ForceSave ]; then
	echo "force-save=${aria2ForceSave}" >> /etc/aria2/aria2.conf
fi
if [ $aria2MaxDownloadLimit ]; then
	echo "max-download-limit=${aria2MaxDownloadLimit}" >> /etc/aria2/aria2.conf
fi
if [ $aria2MaxOverallDownloadLimit ]; then
	echo "max-overall-download-limit=${aria2MaxOverallDownloadLimit}" >> /etc/aria2/aria2.conf
fi
if [ $aria2Continue ]; then
	echo "continue=${aria2Continue}" >> /etc/aria2/aria2.conf
fi
if [ $aria2MinSplitSize ]; then
	echo "min-split-size=${aria2MinSplitSize}" >> /etc/aria2/aria2.conf
fi
if [ $aria2Split ]; then
	echo "split=${aria2Split}" >> /etc/aria2/aria2.conf
fi
if [ $aria2MaxConnectionPerServer ]; then
	echo "max-connection-per-server=${aria2MaxConnectionPerServer}" >> /etc/aria2/aria2.conf
fi
if [ $aria2BtRequireCrypto ]; then
	echo "bt-require-crypto=${aria2BtRequireCrypto}" >> /etc/aria2/aria2.conf
fi
if [ $aria2BtMinCryptoLevel ]; then
	echo "bt-min-crypto-level=${aria2BtMinCryptoLevel}" >> /etc/aria2/aria2.conf
fi
if [ $aria2BtSaveMetadata ]; then
	echo "bt-save-metadata=${aria2BtSaveMetadata}" >> /etc/aria2/aria2.conf
fi
if [ $aria2BtLoadSavedMetadata ]; then
	echo "bt-load-saved-metadata=${aria2BtLoadSavedMetadata}" >> /etc/aria2/aria2.conf
fi
if [ $aria2SaveSessionInterval ]; then
	echo "save-session-interval=${aria2SaveSessionInterval}" >> /etc/aria2/aria2.conf
fi
if [ $aria2ListenPort ]; then
	echo "listen-port=${aria2ListenPort}" >> /etc/aria2/aria2.conf
fi
if [ $aria2EnableDht ]; then
	echo "enable-dht=${aria2EnableDht}" >> /etc/aria2/aria2.conf
fi
if [ $aria2DhtListenPort ]; then
	echo "dht-listen-port=${aria2DhtListenPort}" >> /etc/aria2/aria2.conf
fi

# Extra aria2 option file
if [ ! -f /etc/aria2/aria2_ext.conf ]; then
	touch /etc/aria2/aria2_ext.conf
fi
cat /etc/aria2/aria2_ext.conf >> /etc/aria2/aria2.conf

# Create hook files, if not exists
if [ ! -f /etc/aria2/on-bt-download-complete.sh ]; then
	cp /app/on-bt-download-complete.sh /etc/aria2/on-bt-download-complete.sh
fi
if [ ! -f /etc/aria2/on-download-complete.sh ]; then
	cp /app/on-download-complete.sh /etc/aria2/on-download-complete.sh
fi
if [ ! -f /etc/aria2/on-download-error.sh ]; then
	cp /app/on-download-error.sh /etc/aria2/on-download-error.sh
fi
if [ ! -f /etc/aria2/on-download-pause.sh ]; then
	cp /app/on-download-pause.sh /etc/aria2/on-download-pause.sh
fi
if [ ! -f /etc/aria2/on-download-start.sh ]; then
	cp /app/on-download-start.sh /etc/aria2/on-download-start.sh
fi
if [ ! -f /etc/aria2/on-download-stop.sh ]; then
	cp /app/on-download-stop.sh /etc/aria2/on-download-stop.sh
fi

if [ ! -f /etc/aria2/aria2.session ]; then
	touch /etc/aria2/aria2.session
fi

cat /etc/aria2/aria2.conf

export HOME="~app"
exec s6-setuidgid app /etc/aria2/start.sh
